# server-conf/cron/crontab
0 * * * * sudo indexer --rotate car_ec
0 3 * * * sudo indexer --rotate car_landings_ec
